# lumina/files/_base/usr/etc/profile.d/development.sh
#!/usr/bin/env bash

export TF_CLI_CONFIG_FILE="$HOME/.terraform.d/config.conf"
