# lumina/files/_base/usr/lib/systemd/system/brew-bundle.service
[Unit]
Description=Install Homebrew packages from system and user Brewfiles
After=network-online.target brew-setup.service brew-update.service
Wants=network-online.target brew-setup.service brew-update.service
ConditionPathExists=/etc/.linuxbrew

[Service]
Type=oneshot
User=linuxbrew
WorkingDirectory=/var/home/linuxbrew
EnvironmentFile=/etc/brew/environment
# Install system-managed packages (read-only, part of the image)
ExecStart=/var/home/linuxbrew/.linuxbrew/bin/brew bundle install --file /usr/share/brew/Brewfile
# Install user-managed packages (customizable)
ExecStart=/var/home/linuxbrew/.linuxbrew/bin/brew bundle install --file /etc/brew/Brewfile

[Install]
WantedBy=multi-user.target
