# lumina/files/_base/usr/lib/systemd/system/brew-cleanup.timer
[Unit]
Description=Weekly cleanup of old Homebrew package versions

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
